# system/LINK/etc/init/evm_watchdog.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
exec /bin/evm_watchdog
