# systemd/manageiq-event_handler@.service
[Unit]
PartOf=manageiq-event_handler.target
[Install]
WantedBy=manageiq-event_handler.target
[Service]
WorkingDirectory=/var/www/miq/vmdb
Environment=BUNDLER_GROUPS=manageiq_default,ui_dependencies
EnvironmentFile=/etc/default/manageiq*.properties
ExecStart=/usr/bin/ruby lib/workers/bin/run_single_worker.rb MiqEventHandler --heartbeat --guid=%i
User=manageiq
Restart=no
Type=notify
Slice=manageiq-event_handler.slice
