# systemd/manageiq-remote_console@.service
[Unit]
PartOf=manageiq-remote_console.target
Wants=httpd.service
[Install]
WantedBy=manageiq-remote_console.target
[Service]
WorkingDirectory=/var/www/miq/vmdb
Environment=BUNDLER_GROUPS=manageiq_default,ui_dependencies
EnvironmentFile=/etc/default/manageiq*.properties
ExecStart=/usr/bin/ruby lib/workers/bin/run_single_worker.rb MiqRemoteConsoleWorker --heartbeat --guid=%i
User=manageiq
Restart=no
Type=notify
Slice=manageiq-remote_console.slice
