# build/systemd_units/manageiq-web_service@.service
[Unit]
PartOf=web_service.target
[Install]
WantedBy=web_service.target
[Service]
WorkingDirectory=/var/www/miq/vmdb
Environment=BUNDLER_GROUPS=manageiq_default,ui_dependencies,graphql_api
ExecStart=/bin/bash -lc 'exec ruby lib/workers/bin/run_single_worker.rb MiqWebServiceWorker --heartbeat --guid=%i'
Restart=no
Type=notify
Slice=manageiq-web_service.slice
